# docker/s6-rc.d/rescue-quiescence-reporter/run
#!/command/with-contenv sh
# shellcheck shell=sh
set -eu

key=/run/hermes-rescue-secrets/hmac-current
[ -f "$key" ] || exit 0
[ "$(id -u)" -eq 0 ] || {
    echo "[rescue-reporter] supervision must start as root before setuid" >&2
    exit 111
}

exec s6-setuidgid hermes-rescue \
    /opt/hermes/.venv/bin/python -m agent.rescue_quiescence_reporter serve \
    --runtime-dir /run/hermes-rescue-reporter \
    --continuity-dir /var/lib/hermes-rescue \
    --current-key "$key" \
    --expected-reporter-uid "$(id -u hermes-rescue)" \
    --expected-hermes-uid "$(id -u hermes)"
